# ongoing/autostart-tst/mmiotjs/mmiotjs.service.d/restart.conf
[Service]
Type=forking
PIDFile=mmiotjs.pid
RemainAfterExit=no
Restart=on-failure
RestartSec=5s
