# ongoing/autostart-tst/mmgpio.service.d/restart.conf
[Service]
Type=forking
PIDFile=/var/run/mmgpio.pid
RemainAfterExit=no
Restart=on-failure
RestartSec=5s
