# files/autostart/mmgpio.service.d/restart.conf
[Service]
Type=forking
PIDFile=mmgpio.pid
RemainAfterExit=no
Restart=always
RestartSec=5s
