# ongoing/autostart-tst/mqttbot/mqttbot.service.d/restart.conf
[Service]
Type=forking
PIDFile=mqttbot.pid
RemainAfterExit=no
Restart=on-failure
RestartSec=5s
